# backend/crontab
# 项目定时任务配置
# 
# 智能调度说明：
# - * * * * * : 分时日月周
# - 每15分钟运行一次调度，检查是否有需要拉取的彩种
# - 每小时运行一次补漏检查
# - 每天凌晨运行一次完整的分析统计
#

# 每15分钟运行一次智能调度，检查并拉取当天开奖的彩种
*/15 * * * * cd /app && node backend/scripts/scheduler.js >> /var/log/cron-scheduler.log 2>&1

# 每小时运行一次补漏检查，确保数据完整
0 * * * * cd /app && node backend/scripts/check-and-fill.js >> /var/log/cron-check-fill.log 2>&1

# 每天凌晨 2:00 更新分析统计数据
0 2 * * * cd /app && node backend/scripts/update-stats.js >> /var/log/cron-update-stats.log 2>&1
